# deploy/backup/forge-backup.service
[Unit]
Description=Forge database and artifact backup
Documentation=https://github.com/vivekmaru/task-tracker/blob/main/docs/deployment.md
After=network-online.target

[Service]
Type=oneshot
User=forge
Group=forge
EnvironmentFile=/etc/forge/forge.env
ExecStart=/usr/local/bin/forge-backup.sh
# Only the backup destination needs to be writable.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/backups/forge
