# deploy/backup/forge-backup.timer
[Unit]
Description=Run the Forge backup daily
Documentation=https://github.com/vivekmaru/task-tracker/blob/main/docs/deployment.md

[Timer]
# Daily at 02:30, with catch-up if the machine was off at the scheduled time.
OnCalendar=*-*-* 02:30:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
